# hypr/systemd/user/elephant.service.d/override.conf
[Service]
Environment=ELEPHANT_PROVIDER_DIR=/home/ibanks/.local/lib/elephant/providers-walker
ExecStart=
ExecStart=/home/ibanks/.local/lib/elephant/bin/elephant
